# configs/linux/etc/systemd/system/cpu-performance.service
[Unit]
Description=Set CPU governor to performance for development workloads
After=multi-user.target
ConditionPathExists=/sys/devices/system/cpu/cpu0/cpufreq/scaling_governor

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'for cpu in /sys/devices/system/cpu/cpu*/cpufreq/scaling_governor; do echo performance > "$cpu" 2>/dev/null || true; done'
ExecStart=/bin/bash -c 'echo "CPU governor set to performance mode"'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target